# systemd/diments-gps-tracker.service
[Unit]
Description=DIMENTS GPS Tracker (Node/Express)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/diments-gps-tracker
EnvironmentFile=/opt/diments-gps-tracker/.env
# Jetson Nano + Node 16/18/…: pas dit aan naar jouw node pad indien nodig
ExecStart=/usr/bin/node /opt/diments-gps-tracker/server.js
Restart=always
RestartSec=3
# Hard limits zodat 'ie niet oneindig open files lekt
LimitNOFILE=65535

# Log to journald
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
